# nemo-mobile-session.target
# nemo-mobile.target - a virtual target for the Mobile UI
#
#

[Unit]
Description=Default
