# nemo-middleware.target
#
# nemo-middleware.target - a virtual target for the Mobile UI
#
#

[Unit]
Description=Default
After=dbus.socket
Requires=dbus.socket

