# services/user@.service.d/nemo.conf
[Service]
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%I/dbus/user_bus_socket
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%I
EnvironmentFile=-/var/lib/environment/nemo/*.conf

# Allow hardware adaptations override the default env
EnvironmentFile=-/etc/sysconfig/nemo-mobile-hw
OOMScoreAdjust=-705
# Remove the original ExecStart
ExecStart=
ExecStart=-/lib/systemd/systemd --user --unit=%I.target
